# config/upstart/wompt_auth.conf
description "wompt_auth"

start on startup
stop on shutdown

script
		cd <%=current_path%>/authentication
    exec sudo -u ubuntu rackup --env production 2>>/var/log/wompt.auth.error.log >>/var/log/wompt.auth.log
end script