# nodejs/tools/iptables.conf
description	"iptables"

start on (local-filesystems and net-device-up IFACE!=lo)

script
	iptables -t nat -A PREROUTING -i eth0 -p tcp --dport 80 -j REDIRECT --to-port 16999
end script

